# chirp2mqtt/rootfs/etc/s6-overlay/s6-rc.d/single/run
#!/usr/bin/env bash
set -x
start_chirp_ha.sh
