# chirp2mqtt/rootfs/etc/s6-overlay/s6-rc.d/chirpha/run
#!/usr/bin/env bash
#set -x
#echo "++ chirpha ++"

cd /
exec python3 -m usr.chirpha.start
