# chirp2mqtt/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/env bash
#/command/execlineb -P
set -x
echo "++ nginx ++"

cp /etc/nginx.conf.gtpl /tmp/nginx.conf
nginx -T -c /tmp/nginx.conf

exec nginx -c /tmp/nginx.conf
